# systemd/user/gitbot.timer
[Unit]
Description=Commit and push daily

[Timer]
OnCalendar=daily

[Install]
WantedBy=timers.target
