# systemd/user/gitbot.service
[Unit]
Description=Commit and push select git repo updates

[Service]
Type=oneshot
ExecStart=/home/liam/bin/bashscripts/gitbot
