# systemd/user/test.service
[Unit]
Description=Just for testing purposes
ConditionACPower=true

[Service]
Type=oneshot
ExecStart=/usr/bin/echo "Hello World"

